# official_writeup/prob02-greatwall/greatwall_handout/app/conf.d/openbox.conf
[program:openbox]
command=openbox
autorestart=true