# official_writeup/prob02-greatwall/greatwall_handout/app/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -viewonly -forever -shared
autorestart=true
